# gluon-legacy-https-stats/files/etc/init.d/gluon-legacy-stats
#!/bin/sh /etc/rc.common

START=95

start () {
	iw dev client0 scan passive
	iw dev client1 scan passive
	sh /usr/sbin/ff-stats-wlanscan > /tmp/wlanscan
	sh /usr/sbin/ff-stats
}
